# init/systemd/rbackupd.service
[Unit]
Description=rbackupd backup daemon

[Service]
Type=dbus
BusName=org.rbackupd.daemon
ExecStart=/usr/bin/rbackupd


